# payment_gateway/crontab
0 * * * * python /app/check_status.py >> /var/log/cron.log 2>&1
